# com.local-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local-agent</string>

  <key>ProgramArguments</key>
  <array>
    <!-- python のフルパス: which python3 で確認 -->
    <string>/usr/bin/python3</string>
    <string>/Users/YOUR_USERNAME/mcp-agent-log/local_agent.py</string>
  </array>

  <!-- ログ出力先 -->
  <key>StandardOutPath</key>
  <string>/Users/YOUR_USERNAME/mcp-agent-log/logs/local_agent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/mcp-agent-log/logs/local_agent_err.log</string>

  <!-- クラッシュ時に自動再起動 -->
  <key>KeepAlive</key>
  <true/>

  <!-- ログイン後に起動 -->
  <key>RunAtLoad</key>
  <true/>

  <!-- 作業ディレクトリ -->
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/mcp-agent-log</string>

  <!-- PATH に claude が含まれるように設定 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
